# lirc/files/irexec.init
#!/bin/sh /etc/rc.common

START=99

BIN=/usr/sbin/irexec
NAME=irexec
RUN=/var/run/lirc
PID=$RUN/irexec.pid

start() {
    mkdir -p $RUN
    start-stop-daemon -S -x $BIN -p $PID -m -b
}

stop() {
    start-stop-daemon -K -n $NAME -p $PID -s INT
    rm -rf $PID
}
